# deploy/congestion-forecast.service
# Systemd service for Oracle Cloud / any Linux VPS
# Copy to /etc/systemd/system/congestion-forecast.service
#
# Installation:
#   sudo cp congestion-forecast.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable congestion-forecast
#   sudo systemctl start congestion-forecast
#
# Check status:
#   sudo systemctl status congestion-forecast
#   journalctl -u congestion-forecast -f

[Unit]
Description=Berlin Traffic Congestion Forecast Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/CongestionAI
Environment="PATH=/home/ubuntu/CongestionAI/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ubuntu/CongestionAI/venv/bin/python -m backend.scheduler --daemon
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=congestion-forecast

# Resource limits (adjust as needed for Oracle Free Tier 1GB RAM)
MemoryMax=900M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
